# rootfs/etc/init.d/miio_mpkg
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=81

start_service() {
  procd_open_instance
  procd_set_param command /usr/bin/miio_mpkg
  procd_set_param respawn 3600 5 0
  procd_close_instance
}
